# restic-config/backup.service
[Unit]
Description="Run restic backup"

[Service]
ExecStart=/etc/restic-config/backup.sh
